# nixos/udev-rules/60-i2c.rules
# vim: filetype=udevrules
SUBSYSTEMS=="i2c-dev", GROUP="wheel", RUN+="@setfacl@ -m g:wheel:rwx $devpath"
