# nixos/udev-rules/60-i2c.rules
# vim: filetype=udevrules
SUBSYSTEMS=="i2c-dev", KERNEL=="i2c-[0-9]*", GROUP="wheel", RUN+="@setfacl@ -m g:wheel:rwx $devnode"
